# services/xkcd.service
[Unit]
Description=XKCD
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/home/root/renews.arm \
    -output /usr/share/remarkable/suspended.png \
    -verbose \
    -cooldown COOLDOWN \
    -url https://xkcd.com \
    -xpath '//div[@id="comic"]/img/@src' \
    -mode center \
    -scale 1.75
Restart=always

[Install]
WantedBy=multi-user.target
